# etc/supervisor/conf.d/020-gluster.conf
[program:glusterfs]
command = glusterd --no-daemon
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
